# contrib/cluster-portal.service
[Unit]
Description=video-manager cluster portal (Caddy 背后的管理 API)
After=network.target caddy.service
Wants=caddy.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /home/kxrdyf/scripts/video-manager/contrib/cluster-portal.py
Environment=PEERS_CONF=/etc/caddy/peers.conf
# 用 unix socket 避免 caddy 反代吃 body (HTTP/1.1 + Connection: close 也有问题)
Environment=LISTEN_SOCKET=/run/cluster-portal.sock
Environment=LISTEN_HOST=
Environment=LISTEN_PORT=
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
