# video-manager.service
[Unit]
Description=Video Compression Manager (Web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kxrdyf
WorkingDirectory=/home/kxrdyf/scripts/video-manager
ExecStart=/usr/bin/python3 /home/kxrdyf/scripts/video-manager/app.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/kxrdyf/scripts/video-manager/logs/stdout.log
StandardError=append:/home/kxrdyf/scripts/video-manager/logs/stderr.log

# 安全加固(适度)
NoNewPrivileges=true
ProtectSystem=full
# 这台 NAS volume1 挂在 /home 下,data/logs 都要写到这里,不能 read-only
ProtectHome=false

[Install]
WantedBy=multi-user.target